# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message

## Get the 'message' key from the user config options.
ip=$(bashio::config 'IP')
port="$(bashio::config 'Port')"
passwd="$(bashio::config 'Password')"

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${ip}"
bashio::log.info "${port}"
bashio::log.info "${passwd}"

## Run your program
exec python /usr/bin/raise3d.py
